# Netgear/wnr612v2/package/dnsmasq/files/dnsmasq.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

start() {
	/usr/sbin/dnsmasq -h -r /tmp/resolv.conf
	/sbin/ddns start
}

stop() {
	killall dnsmasq
}
